# packages/gmailnudge/gmailnudge-2.1.1-py3-none-any.whl/gmailnudge/deployment_files/gmailnudge.service
[Unit]
Description=Nudge GMail to force frequent new email checks from my server
After=network-online.target

[Service]
ExecStart=/mnt/share/dev/scripts/gmailnudge/gmailnudge --service

[Install]
WantedBy=multi-user.target
